# ops/systemd/hourkey-job-sweeper.service
[Unit]
Description=Reconcile stale Hourkey durable jobs and refunds
After=docker.service hourkey-decode.service
Requires=docker.service

[Service]
Type=oneshot
WorkingDirectory=/root/releases/current
Environment=NODE_ENV=production
EnvironmentFile=-/etc/hourkey/hourkey.env
ExecStart=/usr/bin/node /root/releases/current/scripts/ops/reconcile-stale-jobs.mjs
TimeoutStartSec=300
