# ops/systemd/hourkey-job-sweeper.timer
[Unit]
Description=Run Hourkey stale job reconciliation hourly

[Timer]
OnBootSec=10m
OnUnitActiveSec=60m
Persistent=true

[Install]
WantedBy=timers.target
